# ictsc2023/gnome-keyring-daemon.conf
[program:gnome-keyring-daemon]
command=dbus-run-session -- gnome-keyring-daemon --unlock
autorestart=true
